# OsQuery Config linux/osquery.conf
{
  "options": {
    "logger_plugin": "filesystem",
    "logger_path": "/var/log/osquery",
    "schedule_splay_percent": "10",
    "disable_logging": "false",
    "utc": "true",
    "verbose": "true"
  },
  "schedule": {
    "Unusual_SSH_Login_T1021_004": {
      "query": "SELECT username, host, time FROM last WHERE tty = 'pts/0' AND (host NOT LIKE '192.168.%' AND host NOT LIKE '10.%');",
      "interval": 300
    },
    "Suspicious_Cron_Job_T1053_003": {
      "query": "SELECT * FROM crontab WHERE command LIKE '%bash%' OR command LIKE '%curl%' OR command LIKE '%wget%';",
      "interval": 300
    },
    "Unauthorized_File_Modification_T1070": {
      "query": "SELECT target_path, action, atime, mtime FROM file_events WHERE target_path LIKE '/etc/%' OR target_path LIKE '/var/log/%';",
      "interval": 300
    },
    "Privilege_Escalation_SUID_T1548_001": {
      "query": "SELECT path, mode FROM file WHERE mode & 4000 != 0;",
      "interval": 300
    },
    "Kernel_Module_Insertion_T1215": {
      "query": "SELECT name, size FROM kernel_modules WHERE name NOT IN ('ext4','xfs','sd_mod','libata','scsi_mod');",
      "interval": 300
    },
    "Hidden_Process_Detection_T1057": {
      "query": "SELECT pid, name, path FROM processes WHERE name NOT LIKE '%(deleted)%' AND pid NOT IN (SELECT pid FROM process_open_sockets);",
      "interval": 300
    },
    "Malicious_Bash_Script_Execution_T1059_004": {
      "query": "SELECT name, path, cmdline FROM processes WHERE name = 'bash' AND (cmdline LIKE '%wget%' OR cmdline LIKE '%curl%' OR cmdline LIKE '%nc %' OR cmdline LIKE '%base64%');",
      "interval": 300
    },
    "Unusual_Network_Connections_T1041": {
      "query": "SELECT local_address, remote_address, remote_port, protocol FROM process_open_sockets WHERE remote_address NOT LIKE '192.168.%' AND remote_address NOT LIKE '10.%';",
      "interval": 300
    },
    "SSH_Authorized_Keys_Abuse_T1098_004": {
      "query": "SELECT target_path, action, mtime FROM file_events WHERE target_path LIKE '%/home/%/.ssh/authorized_keys';",
      "interval": 300
    },
    "Suspicious_Package_Installation_T1092": {
      "query": "SELECT name, version, source FROM rpm_packages WHERE source NOT LIKE '%centos%' AND name NOT LIKE 'gpg-pubkey';",
      "interval": 300
    },
    "Service_Configuration_Change_T1031": {
      "query": "SELECT path, mtime FROM file WHERE path LIKE '/etc/init.d/%' AND size > 0;",
      "interval": 300
    },
    "etc_Passwd_Integrity_Change_T1547_004": {
      "query": "SELECT path, mtime FROM file WHERE path = '/etc/passwd';",
      "interval": 300
    },
    "etc_Hosts_File_Manipulation_T1562": {
      "query": "SELECT path, mtime FROM file WHERE path = '/etc/hosts';",
      "interval": 300
    },
    "DNS_Resolution_Abuse_T1568": {
      "query": "SELECT key, value FROM etc_hosts WHERE value NOT LIKE '127.0.0.1' AND key NOT LIKE 'localhost';",
      "interval": 300
    },
    "Anomalous_Socket_Binding_T1043": {
      "query": "SELECT pid, local_address, local_port FROM listening_ports WHERE local_port NOT IN ('22','80','443','25','53');",
      "interval": 300
    },
    "Process_Injection_Attempt_T1055": {
      "query": "SELECT p1.pid, p1.parent, p1.name, p2.name AS parent_name FROM processes p1 JOIN processes p2 ON p1.parent = p2.pid WHERE p2.name NOT LIKE '%(init|sshd|bash|systemd)%';",
      "interval": 300
    },
    "Abnormal_User_Activity_T1078": {
      "query": "SELECT username, host, time FROM last WHERE time > strftime('%s','now') - 300 AND username NOT LIKE 'root';",
      "interval": 300
    },
    "Unexpected_SetUID_Binaries_T1548_001": {
      "query": "SELECT target_path, action, mtime FROM file_events WHERE target_path LIKE '/usr/bin/%' AND (permissions & 4000) != 0;",
      "interval": 300
    },
    "Bash_History_Extraction_T1552_003": {
      "query": "SELECT target_path, action, mtime FROM file_events WHERE target_path LIKE '%/.bash_history';",
      "interval": 300
    },
    "Insecure_Permission_on_Important_Directories_T1068": {
      "query": "SELECT path, mode FROM file WHERE path IN ('/etc','/usr/local/bin') AND (mode & 2) != 0;",
      "interval": 300
    },
    "Unpatched_Vulnerable_Packages_T1518": {
      "query": "SELECT name, version, release, source, arch FROM rpm_packages WHERE (name LIKE 'openssl' OR name LIKE 'bash' OR name LIKE 'kernel' OR name LIKE 'log4j' OR name LIKE 'log4j-core') AND ((name LIKE 'openssl' AND version < '1.1.1') OR (name LIKE 'bash' AND version < '4.18.0') OR (name LIKE 'kernel' AND version < '5.0.0') OR (name LIKE 'log4j' AND version < '2.17.1') OR (name LIKE 'log4j-core' AND version < '2.17.1'));",
      "interval": 300
    },

    "Ghost_User_Logins": {
      "query": "SELECT l.*, strftime('%Y-%m-%d %H:%M:%S', l.time, 'unixepoch') AS login_time FROM last l LEFT JOIN users u ON l.username = u.username WHERE l.time > strftime('%s','now') - 3600 AND u.username IS NULL;",
      "interval": 60
    }
  },
  "decorators": {
    "always": [
      "SELECT uuid AS host_uuid FROM system_info;"
    ]
  }
}
